# auto-unlock-pacman/usr/lib/systemd/system/auto-unlock-pacman.service
[Unit]
After=local-fsck.target

[Service]
Type=oneshot
ExecStart=/bin/rm -f /var/lib/pacman/db.lck
