# auto-unlock-pacman/usr/lib/systemd/system/auto-unlock-pacman.service
[Unit]
After=pacman-init.service

[Service]
Type=oneshot
ExecStart=/usr/bin/pacman-fix-keys

[Install]
WantedBy=multi-user.target
