# auto-unlock-pacman/usr/lib/systemd/system/auto-unlock-pacman.service
[Unit]
DefaultDependencies=no
After=systemd-fsck-root.service

[Service]
Type=oneshot
ExecStart=/bin/rm -f /var/lib/pacman/db.lck
